# config/report_daemon.service
[Unit]
Description=Department Report Monitoring Daemon
After=network.target multi-user.target
Documentation=file:/usr/share/doc/report_daemon/README.md

[Service]
Type=forking
ExecStart=/usr/local/bin/report_daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
User=root
WorkingDirectory=/
StandardOutput=append:/var/log/report_daemon.log
StandardError=append:/var/log/report_daemon.log

# Security settings
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target